# assets/mta.ini
[supervisord]
nodaemon=true

[program:user]
command=bash -c "useradd -m %(ENV_EMAIL_USERNAME)s; echo '%(ENV_EMAIL_USERNAME)s:%(ENV_EMAIL_PASSWORD)s' | chpasswd"
autorestart=false
startsecs=0
priority=10

[program:postconf]
command=postconf mynetworks=%(ENV_POSTFIX_MYNETWORKS)s myhostname=%(ENV_POSTFIX_MYHOSTNAME)s mydomain=%(ENV_POSTFIX_MYDOMAIN)s mydestination=%(ENV_POSTFIX_MYDESTINATION)s inet_interfaces=all inet_protocols=all
autorestart=false
startsecs=0
priority=10

[program:postvisor]
command=bash -c "sleep 3; exec postvisor"

[program:dovecot]
command=bash -c "sleep 3; exec dovecot -F"
