# assets/mta.ini
[supervisord]
nodaemon=true

[program:dovecot]
command=bash -c "useradd -m %(ENV_EMAIL_USERNAME)s; echo '%(ENV_EMAIL_USERNAME)s:%(ENV_EMAIL_PASSWORD)s' | chpasswd; exec dovecot -F"
priority=10

[program:postvisor]
command=bash -c "sleep 5; exec postvisor"
