# assets/mta.ini
[supervisord]
nodaemon=true

[program:user]
command=bash -c "useradd -m %(ENV_EMAIL_USERNAME)s; echo '%(ENV_EMAIL_USERNAME)s:%(ENV_EMAIL_PASSWORD)s' | chpasswd"
autorestart=false
startsecs=0
priority=10

[program:postvisor]
command=bash -c "sleep 3; exec postvisor"

[program:dovecot]
command=bash -c "sleep 3; exec dovecot -F"
